# l02_stage/raven-bot/services/raven-weight.service
[Unit]
Description=Raven Weight Capture — Flask Web App (bot-iot-l02)
Documentation=https://bot2.sysmayal.ngrok.io
# TODO(Q2): confirm/provision the bot2 ngrok tunnel before go-live
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/raven-bot
EnvironmentFile=/home/admin/raven-bot/rpi_client/.env
Environment=PYTHONUNBUFFERED=1
# NOTE: SCALE_BACKEND/SENSOR_SKILL_ID are reading-tier; web_app.py (this ExecStart) ignores them.
Environment=SCALE_BACKEND=sensor_skill
Environment=SENSOR_SKILL_ID=scale_plant
# NOTE: raven-env must be REBUILT on L02 (venvs are not portable). L01==L02 OS (trixie/py3.13),
# so `python3 -m venv raven-env && raven-env/bin/pip install -r rpi_client/requirements.txt` reproduces it.
ExecStart=/home/admin/raven-bot/raven-env/bin/python3 /home/admin/raven-bot/rpi_client/web_app.py
Restart=always
RestartSec=5
# NO WatchdogSec — web_app.py does not implement sd_notify keepalives, so a
# WatchdogSec here makes systemd SIGABRT-kill the healthy app every interval
# (L01 hit this: ~1295 watchdog-kill restarts, scale dropping offline every ~120s,
# fixed 2026-06-08). Restart=always covers real crashes. Add WatchdogSec back ONLY
# if web_app.py is taught to send WATCHDOG=1 via sd_notify.

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/admin/raven-bot/rpi_client
ProtectHome=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raven-weight

[Install]
WantedBy=multi-user.target
